# sway/.local/service/wayland-pipewire-idle-inhibit/run
#!/bin/ash
exec /usr/bin/wayland-pipewire-idle-inhibit
