# services/.config/rc/init.d/mako
#!/sbin/openrc-run

description="A lightweight notification daemon for Wayland"

command="/usr/bin/mako"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"