# services/.config/rc/init.d/swaylock-suspend
#!/sbin/openrc-run

command="/usr/bin/listen4sleep /usr/bin/swaylockd"
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

