# services/.config/rc/init.d/swaylock-suspend
#!/sbin/openrc-run

command="/home/struan/.local/bin/swaylock-suspend"
supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

