# services/.config/rc/init.d/thunar-daemon
#!/sbin/openrc-run

description="Thunar Daemon"

command="/usr/bin/thunar"
command_args="--daemon"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"
