# sway/.local/service/syncthing/run
#!/bin/ash
exec /usr/bin/syncthing --no-browser 
