# sway/.local/service/thunar-daemon/run
#!/bin/ash
exec /usr/bin/thunar --daemon
