# sway/.local/service/swaylock-suspend/run
#!/bin/ash
exec ~/.local/bin/swaylock-suspend 
