# sway/.local/service/wlsunset/run
#!/bin/ash
exec /usr/bin/wlsunset -l 56.5 -L -3.0
